# python/linux-system-metrics/systemd/metrics-collector.service
[Unit]
Description=Linux VM System Metrics Collector
After=network.target

[Service]
User=deploy
WorkingDirectory=/opt/metrics-collector
ExecStart=/opt/metrics-collector/python-wrapper.sh
Restart=always
RestartSec=5
Environment="PATH=/opt/metrics-collector/venv/bin"

[Install]
WantedBy=multi-user.target
