# python/linux-system-metrics/systemd/metrics-collector.service
[Unit]
Description=Linux VM System Metrics Collector
After=network.target

[Service]
User=deploy
WorkingDirectory=/opt/linux-system-metrics
ExecStart=/usr/bin/python3 /opt/linux-system-metrics/metrics_collector.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
